# trunk/eyerelaxed/Utils/crontab
*/20 * * * * env DISPLAY=:0 ~/runscreen.sh

# m h  dom mon dow   command
